# distsrv.service
[Unit]
Description=App Distribution Service (distsrv)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=distsrv
Group=distsrv
ExecStart=/usr/local/bin/distsrv -config /etc/distsrv/config.toml
Restart=on-failure
RestartSec=5s
WorkingDirectory=/var/lib/distsrv

# Bind privileged ports 80/443 without running as root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/distsrv
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
PrivateDevices=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native

# Memory limits (server has only 1GB RAM)
MemoryMax=300M
MemoryHigh=200M

# Logs go to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
